# backend/app.conf
[program:app]
command=yarn start
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/app.err.log
stdout_logfile=/var/log/supervisor/app.out.log